# deploy/systemd/mdrj-scenario1.service
[Unit]
Description=MDRJ Scenario 1 (centralized logging baseline) - A1
Documentation=https://github.com/alphaserver777/hashgraph
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/mdrj
Environment="PYTHONUNBUFFERED=1"
EnvironmentFile=-/etc/mdrj/scenario1.env
ExecStart=/opt/mdrj/.venv/bin/python -m mdrj.cli node --config /etc/mdrj/scenario1.yaml
Restart=always
RestartSec=5
TimeoutStopSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mdrj-scenario1

# Resource limits — small system services, not production SIEM
LimitNOFILE=4096
MemoryMax=512M

[Install]
WantedBy=multi-user.target
