# deploy/systemd/mdrj.service
[Unit]
Description=MDRJ-DAG distributed peer-to-peer security event registry
Documentation=https://github.com/alphaserver777/hashgraph
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/mdrj
Environment="PYTHONUNBUFFERED=1"
EnvironmentFile=-/etc/mdrj/node.env
ExecStart=/opt/mdrj/.venv/bin/python -m mdrj.cli node --config /etc/mdrj/node.yaml
Restart=always
RestartSec=5
TimeoutStopSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mdrj

LimitNOFILE=4096
# Hashgraph-консенсус требует держать весь видимый DAG в памяти для
# пересчёта famous-witness голосования (O(R²·M²) на recompute). На
# 4 пирах и N~3000 событий пик ~600M, поэтому лимит с запасом 3×.
# Если узел слабый — поднять retention в node.yaml.tpl, чтобы N не рос.
MemoryMax=2G

[Install]
WantedBy=multi-user.target
